# docker/php/cron.d/cron
* * * * * cd /srv/symfony && php bin/console scheduler:consume >> /dev/null 2>&1
